# files/etc/profile.d/homepath.sh
PATH=$HOME/bin:$HOME/.local/share/cargo/bin:$PATH
